# dos-bbs/supervisor/telnetbbs.conf
[program:telnetbbs]
command=/dos/TelnetBBS-master/telnetbbs.pl
directory=/dos/TelnetBBS-master
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
