# dosbox/supervisor/xterm.conf
[program:xterm]
command=/usr/bin/xterm
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
