# dos-bbs-balance/supervisor/balance.conf
[program:balance]
command=/usr/local/bin/startbalance
directory=/dos
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
