# dos-bbs-balance/supervisor/node3.conf
[program:node3]
command=/usr/local/bin/startnode 3
directory=/dos
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
