# dos-bbs-balance/supervisor/node6.conf
[program:node6]
command=/usr/local/bin/startnode 6
directory=/dos
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
