# dos-bbs-balance/supervisor/node2.conf
[program:node2]
command=/usr/local/bin/startnode 2
directory=/dos
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
