# dos-bbs-balance/supervisor/node8.conf
[program:node8]
command=/usr/local/bin/startnode 8
directory=/dos
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
