# dosbox/vncserver.conf
[program:vncserver]
command=/usr/local/bin/startvnc
autostart=true
autorestart=true
startsecs=5

